# scripts/systemd/coredns.service
[Unit]
Description=CoreDNS DNS Server
Documentation=https://coredns.io
After=network.target postgresql.service

[Service]
Type=simple
ExecStart=/usr/local/bin/coredns -conf /etc/dns-stack/Corefile
Restart=on-failure
RestartSec=5
User=root
Group=root

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
